# longview-go.service
[Unit]
Description=Longview client, GO version

[Service]
ExecStart=/root/longview-go.out

[Install]
WantedBy=multi-user.target
